# daemon/upstart/testjob.conf
# /etc/init/testjob.conf

description "A test job file for experimenting with Upstart"
author "Aaron Lelevier"

start on runlevel [2345]

exec echo Test Job ran at  `date` >> /var/log/testjob.log